# docker/cmain/files/upload/crontab
10 4 * * * /usr/sbin/ntpdate time.windows.com
